# systemd/opencpo-proxy.service
[Unit]
Description=OpenCPO Bastion — OCPP WebSocket Proxy
Documentation=https://github.com/opencpo/opencpo-bastion
After=network-online.target opencpo-keyvault.service tailscaled.service
Wants=network-online.target
Requires=opencpo-keyvault.service

[Service]
Type=simple
User=opencpo
Group=opencpo
WorkingDirectory=/opt/opencpo-bastion
ExecStart=/opt/opencpo-venv/bin/python -m gateway.main
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencpo-proxy

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/opencpo /etc/opencpo /tmp

# Resource limits (Pi Zero 2W: 512MB RAM)
MemoryMax=256M
TasksMax=64

[Install]
WantedBy=multi-user.target
