# systemd/opencpo-tap.service
[Unit]
Description=OpenCPO Gateway — OCPP Message Tap (SSE/Query/Export)
After=opencpo-proxy.service

[Service]
Type=simple
User=opencpo
Group=opencpo
WorkingDirectory=/opt/opencpo-gateway
ExecStart=/opt/opencpo-venv/bin/uvicorn gateway.tap:app \
    --host 127.0.0.1 \
    --port 8085 \
    --log-level warning \
    --no-access-log
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencpo-tap

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes

MemoryMax=64M

[Install]
WantedBy=multi-user.target
